# mount-binds
#!/sbin/openrc-run

depend() {
  need localmount
}

start() {
  for MOUNT_POINT in $(find_bind_mounts); do
    do_mount $MOUNT_POINT
  done
}

stop() {
  for MOUNT_POINT in $(find_bind_mounts); do
    do_umount $MOUNT_POINT
  done
}

do_mount() {
    if test -z "$(mount | grep -E "\s$1\s")"; then
      echo "  mount $1"
      mount $1
    else
      echo "  already mounted: $1"
    fi
}

do_umount() {
    if test -z "$(mount | grep -E "\s$1\s")"; then
      echo "  already unmounted: $1"
    else
      echo "  umount $1"
      umount $1
    fi
}

find_bind_mounts() {
  while read LINE; do
    OPTS="$(echo $LINE | awk '{print $4}' | grep -E '(,)?bind(,)?')"
    if test -z "$OPTS"; then
      continue
    fi
    echo $LINE | awk '{print $2}'
  done < /etc/fstab
}

