# open-rc/radarr
#!/sbin/openrc-run
# Copyright 1999-2017 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

DAEMON=/home/daf/apps/Radarr/Radarr.exe
PIDFILE=/var/run/radarr.pid
USER=daf

depend() {
	need net
	need localmount
}

start() {
	ebegin "Starting Radarr"
	chmod 755 $DAEMON
	start-stop-daemon --background --start --exec $DAEMON --make-pidfile --pidfile $PIDFILE --user $USER
	eend $?
}

stop() {
	ebegin "Stopping Radarr"
	start-stop-daemon --stop --pidfile $PIDFILE
	eend $?
}
